# ssh/lib/opt/local/supervisor/conf.d/dropbear.conf
[program:ssh]
command=dropbear -r /opt/local/ssh/dropbear_rsa_host_key  -FEws -p 10022 -P /tmp/ssh/dropbear.pid
stdout_logfile=/opt/local/supervisor/log/ssh-out.log
stdout_logfile_maxbytes=1MB
stderr_logfile=/opt/local/supervisor/log/ssh-err.log
stderr_logfile_maxbytes=1MB
