# base/tags/ubuntu-18.04/opt/app/scripts/crontab
#!/bin/bash
exec fakeroot /usr/bin/crontab -c ${LOCAL}/cron/crontabs -u user "$@"
