# base/tags/alpine-3.8/opt/local/cron/periodic/daily/logrotate
#!/bin/sh
exec /usr/sbin/logrotate -s /opt/local/logrotate/logrotate.status /opt/local/logrotate/logrotate.conf
