# php7/lib/opt/local/supervisor/conf.d/apache2.conf
[program:apache2]
command=/opt/app/apache2/bin/httpd -DFOREGROUND -d /opt/local/apache2 -f /opt/local/apache2/httpd.conf
stdout_logfile=/opt/local/supervisor/log/apache2-out.log
stdout_logfile_maxbytes=1MB
stderr_logfile=/opt/local/supervisor/log/apache2-err.log
stderr_logfile_maxbytes=1MB
