# kafka/lib/opt/local/supervisor/conf.d/kafka.conf
[program:zookeeper]
directory=/opt/local/kafka
command=/opt/local/kafka/bin/zookeeper-server-start.sh ./config/zookeeper.properties
stdout_logfile=/opt/local/supervisor/log/zookeeper-out.log
stdout_logfile_maxbytes=1MB
stderr_logfile=/opt/local/supervisor/log/zookeeper-err.log
stderr_logfile_maxbytes=1MB

[program:kafka]
priority=1000
directory=/opt/local/kafka
command=/opt/local/kafka/bin/kafka-server-start.sh config/server.properties
stdout_logfile=/opt/local/supervisor/log/kafka-out.log
stdout_logfile_maxbytes=1MB
stderr_logfile=/opt/local/supervisor/log/kafka-err.log
stderr_logfile_maxbytes=1MB
