# php7/tags/ubuntu-7.2/opt/local/cron/periodic/15min/clean-php72-session
#!/bin/sh
exec /usr/bin/find /dev/shm/php72/session/ -mindepth 1 -type f -mmin +60 -delete
