# php7/tags/alpine-7.1/opt/local/supervisor/conf.d/php-fpm71.conf
[program:php-fpm71]
command=/opt/app/php71/bin/php-fpm71 -F -p /opt/local/php71 -y /opt/local/php71/php-fpm.conf -c /opt/local/php71
environment=PHP_INI_SCAN_DIR="/opt/local/php71/conf.d"
stdout_logfile=/opt/local/supervisor/log/php-fpm71-out.log
stdout_logfile_maxbytes=1MB
stderr_logfile=/opt/local/supervisor/log/php-fpm71-err.log
stderr_logfile_maxbytes=1MB
